# docs/mitch.service
[Unit]
Description=Mitch Discord Bot - Gaming Buddy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=plexor
Group=plexor
WorkingDirectory=/home/plexor/git/mitch-discord-bot

# Use virtual environment Python
ExecStart=/home/plexor/git/mitch-discord-bot/venv/bin/python3 src/bot.py

# Restart policy
Restart=always
RestartSec=10

# Resource limits
CPUQuota=25%
MemoryMax=2G
MemoryHigh=1.8G

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mitch

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
